# src/oobe/jetson_backend/systemd/oobe-setup.service
[Unit]
# =============================================
# OOBE WiFi Setup Service cho Jetson Nano/Orin Nano
# =============================================
# Service này tự động khởi động khi boot và chạy
# script BLE WiFi Setup để cấu hình mạng lần đầu.
#
# Tương thích:
# - Jetson Nano (JetPack 4.x)
# - Jetson Orin Nano (JetPack 5.x)
# - Jetson Xavier series

Description=OOBE WiFi Setup via BLE for Jetson AI Kit
Documentation=https://github.com/your-repo/jetson-ai-kit

# Đợi các service cần thiết khởi động xong
# - bluetooth.target: BlueZ daemon
# - network-online.target: Network manager
After=bluetooth.target network-online.target
Wants=bluetooth.target

[Service]
# Loại service: forking vì script chạy và giữ
Type=simple

# User chạy service (cần root để truy cập Bluetooth và nmcli)
User=root
Group=root

# Đường dẫn đến script Python
# Lưu ý: Thay đổi đường dẫn này theo vị trí thực tế của script
ExecStartPre=+/bin/systemctl restart bluetooth
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 /opt/oobe-setup/ble_wifi_setup.py

# Khởi động lại nếu script crash
Restart=on-failure
RestartSec=5

# Timeout khởi động (30 giây)
TimeoutStartSec=30

# Tắt systemd watchdog — service dùng idle watchdog nội bộ (10 phút)
WatchdogSec=0

# Hard reset bluetooth sau khi service dừng/crash để clear stale BlueZ state
ExecStopPost=/bin/bash -c 'hciconfig hci0 reset 2>/dev/null; sleep 1; bluetoothctl power on 2>/dev/null'

# Environment variables
Environment=PYTHONUNBUFFERED=1
# Lưu ý: Không cần DBUS_SESSION_BUS_ADDRESS vì script sử dụng SystemBus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oobe-setup

# Working directory
WorkingDirectory=/opt/oobe-setup

[Install]
# Kích hoạt service khi boot vào multi-user mode
WantedBy=multi-user.target
