# setup-firstboot/services/ai-core.service
[Unit]
Description=AI Core Detection Pipeline
After=network.target camera-stream.service
Wants=camera-stream.service

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=__AI_CORE_DIR__
ExecStart=/data/venv/mini-pc/bin/python main.py --device cuda

Restart=always
RestartSec=10

Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY_ENABLED=false

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-core

[Install]
WantedBy=multi-user.target
