# setup-firstboot/services/ai-core.service
[Unit]
Description=AI Core Detection Pipeline
After=network.target camera-stream.service
Wants=camera-stream.service

[Service]
Type=simple
User=__USER__
Group=__USER__

WorkingDirectory=__AI_CORE_DIR__
ExecStart=/usr/bin/python3 main.py --device cuda

Restart=always
RestartSec=10

Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY_ENABLED=false

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-core

[Install]
WantedBy=multi-user.target
