# setup-firstboot/services/backchannel.service
[Unit]
Description=Audio Backchannel Server for go2rtc
After=network.target sound.target pulseaudio.service camera-stream.service
# Restart backchannel when camera-stream restarts (camera-stream recreates echocancel_sink)
BindsTo=camera-stream.service

[Service]
Type=simple
User=__USER__
Group=audio
WorkingDirectory=/opt/backchannel
ExecStart=/bin/bash /opt/backchannel/start.sh
Restart=always
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
