# setup-firstboot/services/backchannel.service
[Unit]
Description=Audio Backchannel Server for go2rtc
After=network.target sound.target pulseaudio.service
# Backchannel only needs echocancel_sink (PulseAudio), independent of go2rtc/stream pipeline.
# go2rtc restart recreates echo cancel in its ExecStartPre; if backchannel's pacat
# loses the old sink → BrokenPipeError → Restart=always → ExecStartPre checks/creates echo cancel.
# Disable rate-limiting so systemd never stops restarting after repeated failures.
StartLimitIntervalSec=0

[Service]
Type=simple
User=__USER__
Group=audio
WorkingDirectory=/opt/backchannel
# Safety net: only create echo cancel if missing. Do NOT run full setup-audio-autostart.sh
# here — it unloads+reloads the module, which would break go2rtc's active pulsesrc connection.
ExecStartPre=/bin/bash -c '\
  U=$(id -u __USER__); \
  export XDG_RUNTIME_DIR=/run/user/$U \
         PULSE_SERVER=unix:/run/user/$U/pulse/native; \
  pactl list short sinks 2>/dev/null | grep -q echocancel_sink \
    || /opt/audio/setup-audio-autostart.sh'
ExecStart=/bin/bash /opt/backchannel/start.sh
Restart=always
RestartSec=5
TimeoutStopSec=5

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
