# setup-firstboot/services/camera-stream.service
[Unit]
Description=Camera Stream Pipeline (CSI + AI + TCP)
After=network.target sound.target pulseaudio.service
Wants=audio-autostart.service

[Service]
Type=simple
User=__USER__
Group=audio
ExecStart=/usr/bin/python3 /opt/stream/start-stream.py --mode service
Restart=always
RestartSec=3
WatchdogSec=60

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
