# setup-firstboot/services/go2rtc.service
[Unit]
Description=go2rtc streaming server (exec: camera pipeline + AI SHM)
After=network.target sound.target pulseaudio.service nvargus-daemon.service

[Service]
Type=simple
User=__USER__
Group=audio

# Reset GPU/camera before start (frees leaked memory from previous runs)
ExecStartPre=+/bin/systemctl restart nvargus-daemon
ExecStartPre=/bin/sleep 2

# Recreate echo cancel BEFORE pulsesrc connects.
ExecStartPre=/bin/bash -c '\
  U=$(id -u __USER__); \
  export XDG_RUNTIME_DIR=/run/user/$U \
         PULSE_SERVER=unix:/run/user/$U/pulse/native; \
  /opt/audio/setup-audio-autostart.sh'

ExecStart=/usr/local/bin/go2rtc -config /etc/go2rtc/go2rtc.yaml

# Re-apply volume AFTER pulsesrc connects.
ExecStartPost=/bin/bash -c '\
  sleep 2; \
  U=$(id -u __USER__); \
  export XDG_RUNTIME_DIR=/run/user/$U \
         PULSE_SERVER=unix:/run/user/$U/pulse/native; \
  pactl set-source-volume echocancel_source 100%% 2>/dev/null; \
  MIC=$(pactl list short sources | grep -iv "monitor\|echo" | grep -i "usb\|jabra" | head -1 | awk "{print \$2}"); \
  [ -n "$MIC" ] && pactl set-source-volume "$MIC" 100%% 2>/dev/null; \
  pactl set-sink-volume echocancel_sink 80%% 2>/dev/null; \
  echo "[audio] Post-connect volume applied: mic=100%% sink=80%%"'

Restart=always
RestartSec=5

# Hard memory limit — kill and restart if leaks beyond 1.5GB
MemoryMax=1536M
MemorySwapMax=512M

# Recycle periodically to prevent long-run memory leaks
RuntimeMaxSec=4h

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
