# setup-firstboot/services/go2rtc.service
[Unit]
Description=go2rtc streaming server
After=network.target camera-stream.service
Wants=camera-stream.service

[Service]
Type=simple
User=__USER__
ExecStart=/usr/local/bin/go2rtc -config /etc/go2rtc/go2rtc.yaml
Restart=always
RestartSec=5

# Hard memory limit — kill and restart if go2rtc leaks beyond 1.5GB
MemoryMax=1536M
MemorySwapMax=512M

# Recycle periodically to prevent long-run memory leaks
RuntimeMaxSec=2h

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
