# setup-firstboot/services/go2rtc.service
[Unit]
Description=go2rtc streaming server
After=network.target camera-stream.service
Wants=camera-stream.service
BindsTo=camera-stream.service

[Service]
Type=simple
User=__USER__
ExecStart=/usr/local/bin/go2rtc -config /etc/go2rtc/go2rtc.yaml
Restart=always
RestartSec=5

# Protect from OOM killer — if memory runs out, kill camera-stream first, not go2rtc
OOMScoreAdjust=-900

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
